# examples/local-linear-codex-bridge.service
[Unit]
Description=Local Linear Codex Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/src/local-linear-codex-bridge
ExecStart=/usr/bin/env bun run src/index.ts daemon --config %h/.config/local-linear-codex-bridge/config.json
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
