# examples/tetherbox.service
[Unit]
Description=Tetherbox
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/path/to/tetherbox
EnvironmentFile=-%h/.config/tetherbox/tetherbox.env
ExecStart=/usr/bin/env bun run src/index.ts daemon --config %h/.config/tetherbox/config.json
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
